# kobrakai.service
[Unit]
Description=KobraKai VoIP Intrusion Prevention System
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/kobrakai/kobrakai.py --config /opt/kobrakai/kobrakai-config.json
WorkingDirectory=/opt/kobrakai
Restart=always
RestartSec=5
KillSignal=SIGINT
KillMode=mixed
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
